# mobile/supabase/ai.portpal.bubble-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.portpal.bubble-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/vbot/.local/share/fnm/aliases/default/bin/node</string>
        <string>bubble-sync.js</string>
        <string>all</string>
        <string>--incremental</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vbot/projects/portpal-web/mobile/supabase</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/vbot/.local/share/fnm/aliases/default/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/vbot</string>
    </dict>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/vbot/Library/Logs/portpal-bubble-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vbot/Library/Logs/portpal-bubble-sync.error.log</string>
</dict>
</plist>
